# config/includes.chroot/etc/cron.d/chef-run
PATH="/sbin:/bin:/usr/local/bin:/usr/local/sbin:/usr/sbin:/usr/bin"
* * * * * root /usr/bin/flock -n /var/tmp/chef-run.lock -c '/usr/bin/chef-run' > /dev/null

