# config/includes.chroot/etc/rc.local
# run chef, then just sit back and let the money roll in...
/usr/bin/flock -n /var/tmp/chef-run.lock -c '/usr/bin/chef-run'

exit 0
